# tasks/xytech-auto-deploy.service
[Unit]
Description=XYTech Website Auto Deploy
After=network.target

[Service]
Type=oneshot
ExecStart=/www/wwwroot/xytech-auto-deploy.sh
User=root
StandardOutput=journal
StandardError=journal
